# scripts/com.clatprep.dailysummary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clatprep.dailysummary</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/arvindkumar/clat_preparation/venv_clat/bin/python</string>
        <string>/Users/arvindkumar/clat_preparation/scripts/daily_summary_job.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/arvindkumar/clat_preparation</string>
    
    <!-- Run at 8:00 AM every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <!-- Also run on load for testing -->
    <key>RunAtLoad</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/arvindkumar/clat_preparation/logs/daily_summary_daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/arvindkumar/clat_preparation/logs/daily_summary_daemon_error.log</string>
</dict>
</plist>

