# scripts/drishti-daily.service
[Unit]
Description=Daily GK pull — Drishti + LegalEdge/TopRankers (PDF + Claude question generation)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/opt/speedmathsgames
ExecStart=/opt/speedmathsgames/venv/bin/python /opt/speedmathsgames/scripts/drishti_daily_cron.py
StandardOutput=append:/opt/speedmathsgames/logs/drishti_daily.log
StandardError=append:/opt/speedmathsgames/logs/drishti_daily.log
# Worst case: 2 sources × 2 dates × 15 min poll = 60 min. 80 min for headroom.
TimeoutStartSec=4800

[Install]
WantedBy=multi-user.target
