# scripts/momentum-scan.service
[Unit]
Description=Daily Momentum Scanner - Minervini Trend Template
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/opt/speedmathsgames
ExecStart=/opt/speedmathsgames/venv/bin/python /opt/speedmathsgames/scripts/momentum_scan_cron.py
Environment=HOME=/opt/speedmathsgames
StandardOutput=append:/opt/speedmathsgames/logs/momentum_scan.log
StandardError=append:/opt/speedmathsgames/logs/momentum_scan.log
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
